# debian/acpi-backlightd.user.service
[Unit]
Description=ACPI backlight Daemon
#Requires=acpid.service

[Service]
ExecStart=/usr/bin/acpi-backlightd

[Install]
WantedBy=default.target
